# packages/net-im/prosody/files/systemd/prosody.service
[Unit]
Description=Jabber/XMPP server
After=network.target

[Service]
Type=forking
User=prosody
Group=prosody
RuntimeDirectory=prosody
RuntimeDirectoryMode=750
PIDFile=/var/run/prosody/prosody.pid
ExecStart=/usr/bin/prosodyctl start
ExecStop=/usr/bin/prosodyctl stop
ExecReload=/usr/bin/prosodyctl reload

[Install]
WantedBy=multi-user.target
